# ndn_python_catalog/ndn-python-catalog.service
[Unit]
Description=NDN Catalog Service

[Service]
Type=simple
ExecStart=/usr/local/bin/ndn-python-catalog

[Install]
WantedBy=multi-user.target
